# templates/supervisor-vortessence-server.ini
{{ ansible_managed | comment }}

[program:vortessence-runserver]
command=python {{ vortessence_rootdir }}/vortessence/manage.py runserver {{ vortessence_server_host }}:{{ vortessence_server_port }}
directory={{ vortessence_rootdir }}/vortessence
user={{ vortessence_user }}
autostart=true
autorestart=false
#autorestart=true
startsecs=10
stopwaitsecs=600
redirect_stderr=true
stdout_logfile=/var/log/supervisor/vortessence-runserver.log

